# Bash_scripts/RabbitMQ_email_consumer_discon/sxrcmonitor.service
[Unit]
Description=SX-registration-collector monitoring service


[Service]
Type=simple
ExecStart=/usr/bin/queue.sh


[Install]
WantedBy=multi-user.target
